# matome/config/supervisord/sc.ini
[program:scraping]
command=/root/.virtualenvs/matome/bin/python /var/flask/matome/matome/manage.py -c ./config/production.py sc
user=root
autorestart=true
stdout_logfile=/var/log/supervisor/scraping-supervisord.log ; 標準出力ログ
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true

